# config/scheduler/crontab
* * * * * /bin/bash /usr/local/bin/magento-cron.sh >> /var/log/cron.log 2>&1
